# scripts/systemd/voice-assistant.service
[Unit]
Description=Voice Assistant (local LLM + STT + TTS)
After=network-online.target ollama.service sound.target
Wants=network-online.target ollama.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_USER__
SupplementaryGroups=audio
WorkingDirectory=__REPO_DIR__
Environment=LOG_LEVEL=INFO
Environment=HOME=/home/__SERVICE_USER__
Environment=PATH=__REPO_DIR__/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=__REPO_DIR__/.venv/bin/uvicorn voice_assistant.main:app \
    --host 0.0.0.0 \
    --port __PORT__ \
    --workers 1 \
    --log-level info
Restart=on-failure
RestartSec=10
MemoryHigh=12G
MemoryMax=14G
TimeoutStartSec=180

StandardOutput=journal
StandardError=journal
SyslogIdentifier=voice-assistant

[Install]
WantedBy=multi-user.target
